# photonOS/4.0GA/files/rc.local
#!/bin/bash

if [ -e /root/ran_customization ]; then
    exit
else
    /root/files/setup.sh &> /var/log/photon-customization.log
fi
